# systemd/rpicamserver.service
[Unit]
Description=Raspberry Pi Camera Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rpicamserver
ExecStart=/usr/bin/env uvicorn app:app --host 0.0.0.0 --port ${APP_PORT:-8000}
Environment=APP_PORT=8000
Restart=on-failure
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
